# debian/fixer-server.service
[Unit]
Description=Fixer aggregation server
Documentation=https://github.com/maumaps/fixer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
EnvironmentFile=-/etc/fixer/fixer-server.env
ExecStart=/usr/bin/fixer-server --config /etc/fixer/fixer-server.toml serve
Restart=on-failure
RestartSec=5
WorkingDirectory=/var/lib/fixer-server
StateDirectory=fixer-server
StateDirectoryMode=0750
AmbientCapabilities=
CapabilityBoundingSet=
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
ProtectSystem=strict
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
